# debian/postinst
#!/bin/bash
echo "Emabling firmware loading at startup"
systemctl enable matrix-creator-firmware
